# focus.service
[Unit]
Description=FOCUS 生活系统 Flask 应用
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/focus-app
# 使用 Gunicorn 生产服务器替代 Flask 内置开发服务器
ExecStart=/home/ubuntu/focus-app/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --access-logfile /home/ubuntu/focus-app/logs/access.log \
    --error-logfile  /home/ubuntu/focus-app/logs/error.log \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
